# templates/etc/systemd/system/nfs-kernel-server.service.d/limits.conf.j2
# {{ ansible_managed }}
[Service]
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
